# supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=10
startretries=1
startsecs=1
